# client/smartcam.service
[Unit]
Description=Smart Camera YOLOv8 WebSocket Client
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/smartcam
ExecStart=/home/pi/smartcam/venv/bin/python /home/pi/smartcam/client.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target


